# site.conf
[program:site]
command=/home/domjudge/dmojsite/bin/uwsgi --ini uwsgi.ini
directory=/home/domjudge/site
stopsignal=QUIT
stdout_logfile=/tmp/site.stdout.log
stderr_logfile=/tmp/site.stderr.log
